# mosquitto.conf
start on net-device-up

respawn

exec /usr/sbin/mosquitto -d -p 80
